# systemd/fan_monitor.service
[Unit]
Description=Thermal Pilot - Smart Fan Monitor Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fan_monitor
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
